# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /var/www/app/deploy/uwsgi.ini